# scripts/rails-research.service
[Unit]
Description=Rails Daily Research Automation
Documentation=https://github.com/jeremedia/rails-8-claude-guide
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rails_daily_research.sh
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/rails-8-claude-guide

# Environment
Environment="RAILS_GUIDE_REPO=/home/deploy/rails-8-claude-guide"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Logging
StandardOutput=append:/var/log/rails-research/rails_research.log
StandardError=append:/var/log/rails-research/rails_research_error.log

# Security
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/deploy/rails-8-claude-guide /var/log/rails-research

# Restart policy
Restart=on-failure
RestartSec=300

[Install]
WantedBy=multi-user.target